# AstroFinSentinelV5/deploy/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:sentinel]
command=python -m orchestration.sentinel_v5
directory=/app
stdout_logfile=/var/log/supervisor/sentinel.log
stderr_logfile=/var/log/supervisor/sentinel_err.log
autostart=true
autorestart=true
environment=PYTHONPATH="/app"

[program:health]
command=python -m uvicorn deploy.monitoring.health_endpoints:app --host 0.0.0.0 --port 8000
directory=/app
stdout_logfile=/var/log/supervisor/health.log
stderr_logfile=/var/log/supervisor/health_err.log
autostart=true
autorestart=true
